# master/checkusbstick.service
[Unit]
Description=Grab new messages files from usb stick

[Service]
Type=simple
ExecStart=/usr/local/bin/checkusbstick.sh
